# server/lead-finder.service
[Unit]
Description=Lead Finder API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/ubuntu/projects/lead-finder-api
EnvironmentFile=/home/ubuntu/projects/lead-finder-api/.env
ExecStart=/usr/bin/node /home/ubuntu/projects/lead-finder-api/dist/src/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/projects/lead-finder-api

[Install]
WantedBy=multi-user.target
