# server/lead-finder.service
[Unit]
Description=Lead Finder API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/lead-finder/server
EnvironmentFile=/opt/lead-finder/server/.env
ExecStart=/usr/bin/node /opt/lead-finder/server/dist/src/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/lead-finder/server

[Install]
WantedBy=multi-user.target

